# site.conf
# this goes in /etc/init/
# to start run 
# sudo initctl start site

start on startup
chdir /home/mark/Euterpea_Web
exec /home/mark/.cabal/bin/EuterpeaSite -p 80
